# templates/join_swarm.service
[Unit]
Description=Ensure that node joins a swarm on startup

[Service]
ExecStart=/root/join_swarm.sh

[Install]
WantedBy=multi-user.target